# adb-addon/rootfs/etc/services.d/adb/run
#!/usr/bin/with-contenv bash
set -e

echo "Starting ADB Server..."
exec adb -a -P 5037 server nodaemon
